# deploy/coffee-ai-boss.service
[Unit]
Description=Coffee AI Boss (FastAPI / uvicorn)
Documentation=https://github.com/flying-dragon-ai/coffee-ai-boss
After=network.target
# 若 MySQL / Redis 也在本机用 systemd 管理，取消下两行注释以强制启动顺序：
# Wants=mysql.service redis-server.service
# After=mysql.service redis-server.service

[Service]
Type=simple
User=coffee
Group=coffee
WorkingDirectory=/opt/coffee-ai-boss
EnvironmentFile=/opt/coffee-ai-boss/.env
ExecStart=/opt/coffee-ai-boss/.venv/bin/python -m uvicorn app.main:app --host 0.0.0.0 --port 8000 --workers 1
Restart=always
RestartSec=3
KillSignal=SIGINT
# 日志走 journald，用 `journalctl -u coffee-ai-boss -f` 查看
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
