# deploy/crossroads-agent-cafe.service
[Unit]
Description=Crossroads Agent Cafe (FastAPI / uvicorn)
Documentation=https://github.com/flying-dragon-ai/crossroads-agent-cafe
After=network.target

[Service]
Type=simple
User=coffee
Group=coffee
WorkingDirectory=/opt/crossroads-agent-cafe
EnvironmentFile=/opt/crossroads-agent-cafe/.env
Environment=ENVIRONMENT=production
Environment=AUTH_COOKIE_SECURE=true
ExecStartPre=/opt/crossroads-agent-cafe/.venv/bin/python /opt/crossroads-agent-cafe/scripts/migrate_order_sources.py
ExecStart=/opt/crossroads-agent-cafe/.venv/bin/python -m uvicorn app.main:app --host 0.0.0.0 --port 8000 --workers 1 --ws-max-size 16384
Restart=always
RestartSec=3
KillSignal=SIGINT
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
